# crontab
0 6 * * * root /opt/hexo/git.sh all >> "/home/cron/cron.log 2>&1

