# crontab
SHELL=/bin/bash
0 6 * * * /opt/hexo/git.sh >> "/home/cron/cron.log" 2>&1 &
# */2 * * * * /opt/hexo/git.sh >> "/home/cron/cron.log" 2>&1 &
